# rust-media-server/rust-media-server.service
[Unit]
Description=Neustream Rust Media Server
After=network.target

[Service]
Type=simple
User=neustream
Group=neustream
WorkingDirectory=/opt/neustream
ExecStart=/opt/neustream/rust-media-server --config /opt/neustream/config.toml
Restart=always
RestartSec=5
Environment=RUST_LOG=info
Environment=CONTROL_PLANE_URL=https://api.neustream.app
Environment=RTMP_PORT=1935
Environment=HTTP_PORT=8000

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/neustream

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target